# systemd/caption-watchdog.timer
[Unit]
Description=Run caption watchdog every minute

[Timer]
OnBootSec=90
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
