# systemd/caption-monitor.service
[Unit]
Description=Gramps Caption Service Monitor

[Service]
Type=oneshot
ExecStart=%h/gramps-transcriber/scripts/caption-monitor.sh
