# systemd/caption-watchdog.service
[Unit]
Description=Caption Service Watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/caption-watchdog.sh
