# systemd/display-watchdog.service
[Unit]
Description=Display Watchdog - Ensures LightDM is running
After=lightdm.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/display-watchdog.sh
