# systemd/network-watchdog.timer
[Unit]
Description=Run network watchdog every 2 minutes

[Timer]
OnBootSec=120
OnUnitActiveSec=120

[Install]
WantedBy=timers.target
